# docker/supervisord-post-deploy.conf
[program:post-deploy]
command=/home/app/bin/post-deploy.sh
priority=1
startsecs=0
autorestart=false
stdout_logfile=syslog
stderr_logfile=syslog
